# radarr.ini
[program:radarr]
user=
command=/usr/bin/mono /opt/Radarr/Radarr.exe /nobrowser /data=/config/
