# radarr.ini
[program:radarr]
user=
command=/usr/bin/mono /opt/radarr/NzbDrone.exe /nobrowser /data=/config/
