# radarr.ini
[program:radarr]
user=
Environment=MONO_TLS_PROVIDER=legacy
command=/usr/bin/mono /opt/Radarr/Radarr.exe /nobrowser /data=/config/
